# docker/etc/service/puma/run
#!/bin/bash

exec 2>&1
cd /web_app

exec puma
